# daemon/com.mysc.yapnotes.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.mysc.yapnotes</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>YAPNOTES_ROOT_PLACEHOLDER/daemon/yapnotesd.sh</string>
    </array>

    <!-- 로그인 시 자동 시작 -->
    <key>RunAtLoad</key>
    <true/>

    <!-- 크래시/종료 시 자동 재시작 -->
    <key>KeepAlive</key>
    <true/>

    <!-- 로그 파일 -->
    <key>StandardOutPath</key>
    <string>LOG_FILE_PLACEHOLDER</string>
    <key>StandardErrorPath</key>
    <string>LOG_FILE_PLACEHOLDER</string>

    <!-- PATH: Homebrew + swift 경로 포함 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>HOME_PLACEHOLDER</string>
    </dict>

    <!-- 재시작 간격 (5초 이내 재시작 반복 시 30초 대기) -->
    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
